# .config/systemd/user/gammastep.service
[Unit]
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/gammastep

[Install]
WantedBy=sway-session.target
