# .config/systemd/user/mako.service
[Unit]
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/usr/bin/mako

[Install]
WantedBy=sway-session.target
