# .config/systemd/user/mbsync.service
[Unit]
Description=mbsync as a service

[Service]
Type=simple
ExecStart=%h/.local/services/mbsync-run
