# .config/systemd/user/sway-idle.service
[Unit]
BindsTo=sway-session.target

[Service]
Type=simple
ExecStart=/home/meril/.local/services/sway-idle

[Install]
WantedBy=sway-session.target
