# .config/systemd/user/goimapnotify@.service.d/override.conf
[Unit]
After=graphical-session.target
PartOf=graphical-session.target

[Install]
WantedBy=
WantedBy=graphical-session.target