# .config/systemd/user/goimapnotify@.service.d/override.conf
[Service]
ExecCondition=/bin/sh -c '[ -n "$WAYLAND_DISPLAY" ]'