# .config/systemd/user/gpg-agent.service.d/override.conf
[Service]
PassEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR XDG_CURRENT_DESKTOP
